# crontab
* */23 * * * python /whoisAutomated.py >> /logfile
